# telnetecho.socket
[Socket]
ListenStream=24
Accept=true

[Install]
WantedBy=sockets.target

